# agents/buildagent/bin/jetbrains.teamcity.BuildAgentUpgrade.plist.dist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>WorkingDirectory</key>
	<string>####WORKDIR</string>
	<key>Debug</key>
	<false/>
	<key>Label</key>
	<string>jetbrains.teamcity.BuildAgentUpgrade</string>
	<key>OnDemand</key>
	<false/>
	<key>KeepAlive</key>
	<true/>
	<key>ProgramArguments</key>
	<array> 
<string>####ARZ</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>logs/launchd.upgrade.err.log</string>
	<key>StandardOutPath</key>
	<string>logs/launchd.upgrade.out.log</string>
</dict>
</plist>
